# src/DEBIAN/postrm
#!/bin/bash
PACKAGE_NAME="$DPKG_MAINTSCRIPT_PACKAGE"

if [ "${PACKAGE_NAME}" = "cn.flamescion.bookworm-compatibility-mode" ];then
systemctl stop ace-bookworm-auto-upgrade
systemctl disable ace-bookworm-auto-upgrade
fi



if [ "$1" = "remove" ] || [ "$1" = "purge" ];then

echo "清理卸载残留"

rm -rf /opt/apps/$PACKAGE_NAME
else
echo "非卸载，跳过清理"
fi



